# salt/base/software/postgresql/files/postgresql.service.jinja
{% set dbconf = salt['pillar.get']('postgresql_conf') -%}
.include /lib/systemd/system/postgresql.service

[Service]
Environment=PGDATA={{ dbconf.homedir }}
