# salt/software/postgresql/files/postgresql.service.jinja
{% set dbconf = salt['pillar.get']('postgresql_conf') -%}
.include /lib/systemd/system/postgresql-{{- dbconf.version -}}.service

[Service]
Environment=PGDATA={{ dbconf.homedir }}/data
